# pipeline/phase_e/robops-camera-bridge.service
[Unit]
Description=RobOps camera bridge — picamera2 → UDP MJPEG → Docker camera_node
After=network.target
# Restart if the Docker stack restarts and frees the camera
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=serkanrob
# Pi OS system-installed picamera2 lives in dist-packages
Environment="PYTHONPATH=/usr/lib/python3/dist-packages"
ExecStart=/usr/bin/python3 /opt/robops/camera_bridge.py \
    --host 127.0.0.1 \
    --port 5000 \
    --width 640 \
    --height 480 \
    --fps 5 \
    --jpeg-quality 75
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
